# actividad4/ej_servicio.service
[Unit]
Description=Ejecucion de script de saludo
[Service]
User=pind20
ExecStart=/home/script_act4/script.sh
[Install]
WantedBy=multi-user.target